# media-processor-ui.service
[Unit]
Description=Media Processor Web UI
After=network.target
Wants=media-processor-py.service media-processor-api.service

[Service]
Type=simple
User=sharvinzlife
WorkingDirectory=/home/sharvinzlife/media-processor/web-app
ExecStart=/usr/bin/node /home/sharvinzlife/media-processor/web-app/server.js
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=media-processor-ui
Environment=NODE_ENV=production
Environment=PORT=3001

[Install]
WantedBy=multi-user.target 