# media-processor-web.service

[Unit]
Description=Media Processor Web Interface
After=network.target

[Service]
WorkingDirectory=/home/sharvinzlife/media-processor/web-app/api
ExecStart=/usr/bin/node server.js
Restart=always
User=sharvinzlife
Environment=PORT=3001
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target